# scripts/phone-keepawake-launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.longxin.zuozhile.keepawake</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/zhanglongxin/Library/Zuozhile/phone-keepawake.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>480</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/zhanglongxin/Library/Logs/zuozhile-phone-keepawake.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/zhanglongxin/Library/Logs/zuozhile-phone-keepawake.launchd.err.log</string>
</dict>
</plist>
